# runit/service/patroni/run
#!/usr/bin/env sh
set -eu

exec 2>&1

COMMAND=/usr/bin/patroni
CONFIG=/etc/patroni/postgres0.yml
USER=postgres

exec /sbin/chpst -u ${USER} ${COMMAND} ${CONFIG}