# runit/service/haproxy/run
#!/usr/bin/env sh
set -eu

exec 2>&1

COMMAND=/usr/sbin/haproxy
CONFIG=/etc/haproxy/haproxy.cfg
USER=haproxy
GROUP=haproxy

exec /sbin/chpst -u ${USER}:${GROUP} ${COMMAND} -- ${CONFIG}